# files/etc/uci-defaults/mi3g
#r3g 分支main 
#开WiFi 5g WiFi信道设置为国家_俄罗斯
uci set wireless.radio0.disabled='0'
uci set wireless.radio0.country='RU'
uci set wireless.radio0.channel='11'
uci set wireless.radio1.disabled='0'
uci set wireless.radio1.channel='64'
uci set wireless.radio1.country='RU'
uci commit wireless

#设置 usbnwan usbwan6 
uci set network.usbwan=interface
uci set network.usbwan.proto='dhcp'
uci set network.usbwan.device='eth1'
uci set network.usbwan6=interface
uci set network.usbwan6.proto='dhcpv6'
uci set network.usbwan6.device='eth1'
uci commit network
                  
#设置usbwan usbwan6 加入防火墙 并开启硬件加速
uci add_list firewall.cfg03dc81.network='wan'
uci add_list firewall.cfg03dc81.network='wan6'
uci add_list firewall.cfg03dc81.network='usbwan'
uci add_list firewall.cfg03dc81.network='usbwan6'
uci set firewall.cfg01e63d.flow_offloading='1'
uci set firewall.cfg01e63d.flow_offloading_hw='1'
uci commit firewall

echo "All done!"
